# 04_rescue_ssh/pgen.service
[Unit]
Description=Password Generator Microservice (pgen)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pgen --serve --port 7777 --host 127.0.0.1
Restart=on-failure
RestartSec=5s

# Security and Sandboxing (systemd-native hardening)
DynamicUser=yes
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
